# systemd/inv1.service
[Unit]
Description=Fetch Data from inverter

[Service]
Type=simple
ExecStart=/usr/local/bin/mqtt-push1.sh

[Install]
WantedBy=multi-user.target
