# systemd/inv1.timer
[Unit]
Description=Update Inverters every minute

[Timer]
OnCalendar=*:*:00/5

[Install]
WantedBy=timers.target
