# systemd/inv2.timer
[Unit]
Description=Update Inverters every minute

[Timer]
OnCalendar=*:00/2:12

[Install]
WantedBy=timers.target
